# _docker/supervisord.conf
[supervisord]
nodaemon = true
logfile = /dev/stdout
loglevel = info

[unix_http_server]
file = /run/supervisord.s

[supervisorctl]
serverurl = unix:///run/supervisord.s

[program:masif-upgrader-agent]
command = /bin/bash -exo pipefail -c 'perl -pi -e "s/(?<=agent)X/$MUAGENT/" /etc/masif-upgrader-agent.ini; exec masif-upgrader-agent -config /etc/masif-upgrader-agent.ini'
depends_on = mariadb, influxdb
stdout_logfile = /dev/stdout
redirect_stderr = true
autostart = true
autorestart = true
startsecs = 10
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 5

[program:mkfakepkg]
command = /bin/bash -exo pipefail -c 'while sleep 30; do mkfakepkg; apt-get update; done'
stdout_logfile = /dev/stdout
redirect_stderr = true
autostart = true
autorestart = true
startsecs = 10
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 5
